# influx-import.service
[Unit]
Description = Importing temperature and humidity data to influxDB
After = network.target 
 
[Service]
Type = simple
ExecStart = python3 /home/pi/influx_import.py 
User = pi
Group = pi
Restart = on-failure 
SyslogIdentifier = influx-data
RestartSec = 5
TimeoutStartSec = infinity
 
[Install]
WantedBy = multi-user.target
